# websocket.conf
[program:laravel-websocket]
command=php artisan websocket
directory=/var/www/codeblock