# worker.conf
[program:laravel-redis]
command=php artisan queue:work redis --tries=3
directory=/var/www/codeblock
#directory=/home/vagrant/Code/Codeblock